# systemd/sould.service
[Unit]
Description=Service for mirroring git repositories
Requires=sould-gitd.service

[Service]
ExecStart=/usr/bin/sould
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
